# config/profile.d/50-dircolors.sh
# This eval does also `export LS_COLORS`.
eval "$(dircolors --sh "${XDG_CONFIG_HOME:-${HOME}/.config}/dir_colors")"
