# config/profile.d/80-user-bin.sh
PATH="${XDG_BIN_HOME:-"$(readlink -f "${XDG_DATA_HOME:-"${HOME}/.local/share"}/../bin")"}:${PATH}"
